# supervisord.conf
[supervisord]
nodaemon=true
[program:jenkins]
command=/etc/init.d/jenkins start
